# templates/stockprice.service.j2
[Unit]

[Service]
ExecStart=/bin/sh -c "sudo mvn spring-boot:run"
WorkingDirectory=/home/{{ node_user }}/stockprice